# roles/iptables-script/files/etc_init.d_iptables-script
#!/sbin/openrc-run

name="iptables-script"

depend() {
    before net
}

start() {
    ebegin "Running /usr/local/sbin/iptables-script"
    /usr/local/sbin/iptables-script
    eend $? "Failed to start ${name}"
}
